# ansible/roles/media-edge-install/files/edge/hnsd/hnsd.service.j2
[Unit]
Description=HNSD - Handshake Daemon

[Service]
ExecStart=/opt/hnsd/hnsd -p 4 -r 127.0.0.1:53
Restart=always
User=root

[Install]
WantedBy=multi-user.target